# systemd/user/update-reminder.service
[Unit]
Description=Weekly Ubuntu update reminder notification

[Service]
Type=oneshot
ExecStart=/home/wei-ying-chen/.local/bin/update-reminder.sh
